# supervisor.conf
[supervisord]
nodaemon=true

[program:startMongo]
command=/bin/bash -c "/bin/bash /home/workweek-user/start-mongo.sh"

[include]
files = /etc/supervisor/conf.d/*.conf
